# deploy/taaip.service.example
[Unit]
Description=TAAIP FastAPI app
After=network.target

[Service]
Type=simple
User=taaip
Group=www-data
WorkingDirectory=/opt/TAAIP
EnvironmentFile=/opt/TAAIP/.env
ExecStartPre=/usr/bin/env zsh /opt/TAAIP/scripts/taaip_preflight.sh
ExecStart=/opt/TAAIP/.venv/bin/python -m uvicorn services.api.app.main:app --host ${HOST:-127.0.0.1} --port ${PORT:-8001}
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target
